# contrib/upstart.conf.example
description "IRCD Load Balancer"
author      "prawnsalad - http://guthub.com/prawnsalad/"
 
start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn

exec /usr/bin/env node /opt/ircdbalancer/ircdbalancer.js -d >> /var/log/ircdbalancer.log 2>&1
